# dot_config/systemd/user/waybar.service
[Unit]
Description=Wayland bar for Hyprland
Partof=graphical-session.target

[Service]
ExecStart=/usr/bin/waybar
Restart=always
RestartSec=2

[Install]
WantedBy=graphical-session.target
